# etc/rc.local
/usr/local/sbin/dhcp6c vmx0
/usr/sbin/rtadvd vmx1
